# deploy/cq-webhook.service
[Unit]
Description=CQ v3 AI deploy webhook
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/root/cq-v3-ai
EnvironmentFile=/etc/cq-webhook.env
Environment=REPO_DIR=/root/cq-v3-ai
Environment=DEPLOY_BRANCH=main
Environment=WEBHOOK_PORT=9000
ExecStart=/usr/bin/python3.11 /root/cq-v3-ai/deploy/webhook.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
